# docker/php/supervisord-job.conf
; job コンテナ用 supervisord 設定
; 通常の Job (queue=default) を処理する queue:work を 1 プロセス常駐
;
; このプロジェクトでは Bus::batch の Job も含め、すべての非同期 Job は
; default queue に流して job コンテナで処理する (queue 分離はしない)。
; 重い処理を別ワーカーに分離したい場合は、必要になった時点で queue 名を
; 切り出して別 supervisord program を追加する。
;
; Laravel 公式推奨パターン (https://laravel.com/docs/queues#supervisor-configuration)

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; -----------------------------------------------------------------------------
; job-worker (queue=default を処理 / Bus::batch の Job もここに来る)
;   --max-time=3600  1 時間で自己終了 → supervisord が再起動
;   --memory=256     256MB 超で自己終了 (リーク対策)
;   stopwaitsecs=60  停止シグナルから kill までの猶予 (実行中ジョブを完走させる)
; -----------------------------------------------------------------------------
[program:job-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=default --tries=3 --backoff=10 --max-time=3600 --memory=256 --sleep=3
user=www-data
autostart=true
autorestart=true
numprocs=1
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
