# docker/php/supervisord-batch.conf
; batch コンテナ用 supervisord 設定
; cron daemon (Laravel 公式推奨パターン) を supervisord 配下で動かす。
; batch = 「時刻トリガで起動する定期処理 (= Linux 一般語の batch)」 という位置づけ。
; Laravel の Bus::batch とは別概念で、こちらの batch コンテナは Job を直接実行しない。
;
; 構成:
;   - [program:cron]      : Debian cron daemon をフォアグラウンドで起動
;   - [program:cron-tail] : cron job の出力を tail -F して stdout に流す
;                           (docker compose logs batch で見える)
;
; なぜ schedule:work / シェルループではなく cron daemon か:
;   - Laravel 公式の Production 推奨が cron + schedule:run
;   - 運用ツール / モニタリングが cron 前提で揃っていることが多い
;   - Linux の標準スケジューラと挙動 (PATH / MAILTO / 環境変数) が一致する
;   - Laravel 10 / 11 / 12 / 13 すべてで同じ書き方

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; -----------------------------------------------------------------------------
; cron daemon (フォアグラウンド起動)
;   -f : フォアグラウンド (supervisord に監視させるため)
;   -L 4 : ログレベル (job 開始/終了をログに出す)
; -----------------------------------------------------------------------------
[program:cron]
command=/usr/sbin/cron -f -L 4
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------------
; cron-tail: scheduler-cron.log を tail -F して supervisord 経由で stdout に流す
;   これで `docker compose logs batch` から schedule:run の出力が見える
;   - --retry: 起動時にファイルがなくても待つ
;   - -F (大文字): rotate 等で inode が変わっても追従
; -----------------------------------------------------------------------------
[program:cron-tail]
command=/bin/sh -c "mkdir -p /var/www/html/storage/logs && touch /var/www/html/storage/logs/scheduler-cron.log && chown www-data:www-data /var/www/html/storage/logs/scheduler-cron.log && exec tail -F /var/www/html/storage/logs/scheduler-cron.log"
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
